# alphagenesisdashboard.service
[Unit]
Description=AlphaGenesis Real-time Dashboard
After=network.target

[Service]
Type=simple
User=outerwoolf
WorkingDirectory=/home/outerwoolf/AlphaGenesis
Environment=PATH=/home/outerwoolf/AlphaGenesis/venv/bin
ExecStart=/home/outerwoolf/AlphaGenesis/venv/bin/python dashboard/simple_dashboard_5001.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 